# psvibe-attendance.service
[Unit]
Description=PS VIBE Staff Attendance System
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/staff_attendance
ExecStart=/usr/bin/node /opt/staff_attendance/staff_attendance.js server 3099
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
