# ibet789-bot/ibet789-bot.service
[Unit]
Description=iBet789 Telegram Bot — Agent Dashboard Automation
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ibet789-bot
ExecStart=/usr/bin/node /opt/ibet789-bot/bot.js
Restart=always
RestartSec=5
User=root
EnvironmentFile=/opt/ibet789-bot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ibet789-bot

# Security hardening (optional)
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
