# refactor_staging/psvibe-bot.service
[Unit]
Description=PS Vibe Staff Bot
After=network.target
[Service]
Type=simple
User=root
WorkingDirectory=/root/Sales-Tele-Bot
EnvironmentFile=/root/Sales-Tele-Bot/.env
ExecStart=/root/Sales-Tele-Bot/.venv/bin/python3 main.py
Restart=always
RestartSec=3
TimeoutStopSec=10
KillSignal=SIGINT
StandardOutput=append:/root/Sales-Tele-Bot/logs/bot.log
StandardError=append:/root/Sales-Tele-Bot/logs/bot.log
[Install]
WantedBy=multi-user.target
